# conf/minion_supervisord.conf
[supervisord]
nodaemon=true

[program:ssh]
command=service ssh start
autorestart=true

[program:wrapdocker]
command=/usr/local/bin/wrapdocker
autorestart=true