# conf/supervisord/00_default.conf
[supervisord]
nodaemon=true

[program:ssh]
command=service ssh start
autorestart=true