# conf/supervisord/01_minion.conf
[supervisord]
nodaemon=true

[program:ssh]
command=service ssh start
autorestart=true

[program:wrapdocker]
command=/bin/wrapdocker
autorestart=false
